# music/to-load.cron
# do daily/weekly/monthly maintenance
# min   hour    day     month   weekday command
24       7       *       *       1       nice beet fetchart
24       7       *       *       *       nice beet splupdate
24       7       *       *       2       date >> /config/automation.log;echo "beet move" >> /config/automation.log;nice beet move >> /config/automation.log; date >> /config/automation.log;
24       7       *       *       3       date >> /config/automation.log;echo "beet update" >> /config/automation.log;nice beet update >> /config/automation.log; date >> /config/automation.log;
24       6       *       *       *       date >> /download/beets-ingest/activity.log; beet import -piqm -l /download/beets-ingest/import.log /download/beets-ingest/ >> /download/beets-ingest/activity.log; beet move >> /config/automation.log; date >> /download/beets-ingest/activity.log;
*/5      *       *       *       *       date > /tmp/cron-running
24       7       *       *       3       cd /config; gzip -k musiclibrary.db
24       7       *       *       5       truncate -s 10M /config/automation.log